# com.budgetmanagement.servers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.budgetmanagement.servers</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/Samsung USB/budget_management_anti/start_all_servers.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Volumes/Samsung USB/budget_management_anti/logs/launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/Samsung USB/budget_management_anti/logs/launchd_error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/Samsung USB/budget_management_anti</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>

